# launchd/com.miya.cloudflared-mac-activity.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miya.cloudflared-mac-activity</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/miya/Blogs/scripts/run_cloudflared_tunnel.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/miya/Blogs</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/miya/Library/Logs/cloudflared-mac-activity.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/miya/Library/Logs/cloudflared-mac-activity.error.log</string>
</dict>
</plist>
